# launchagents-template/com.example.mc-mailbox-digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Template: MC Assistant Mailbox daily digest. Fires daily 21:00.
  Wrapper: scripts/cron-mc-mailbox-digest.sh — silent if all counters are 0.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.example.mc-mailbox-digest</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/YOU/chief-of-staff/scripts/cron-mc-mailbox-digest.sh</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/YOU/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/YOU</string>
    </dict>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>21</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/YOU/chief-of-staff/logs/cron-mc-mailbox-digest-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/YOU/chief-of-staff/logs/cron-mc-mailbox-digest-stderr.log</string>
</dict>
</plist>
